# templates/backup.service.j2
[Unit]
Description=Backup system wit Restic

[Service]
Type=simple
ExecStart=/bin/bash {{ duplicity_download_path }}restic/{{utility.lower()}}_backup.sh 

[Install]
WantedBy=multi-user.target
